# data/examples/35/redis.container
[Container]
Image=redis:3.0
Pod=example.pod
