# data/examples/35/web.container
[Container]
Image=tuna/docker-counter23
Pod=example.pod
