# data/fixtures/example-proxy/deploy/example-proxy/example-proxy_postgres.container
[Container]
Image=docker.io/library/postgres:17-alpine
Volume=example-proxy_postgres-data.volume:/var/lib/postgresql/data:Z
Volume=./initdb:/docker-entrypoint-initdb.d:ro,Z
Environment=POSTGRES_USER=app
Environment=POSTGRES_PASSWORD_FILE=/run/secrets/db_password
Environment=POSTGRES_DB=myapp
Environment=PGDATA=/var/lib/postgresql/data/pgdata
Environment=TZ=Australia/Melbourne
Network=example-proxy_backend.network
ContainerName=postgres
DropCapability=ALL
AddCapability=DAC_OVERRIDE
AddCapability=SETUID
AddCapability=SETGID
AddCapability=FOWNER
AddCapability=CHOWN
ReadOnly=true
Tmpfs=/tmp:noexec,nosuid,size=256m
Tmpfs=/run/postgresql:size=64m
ShmSize=256m
Sysctl=net.core.somaxconn=1024
LogDriver=json-file
HealthCmd=pg_isready -U app -d myapp
HealthInterval=10s
HealthRetries=5
HealthTimeout=5s
HealthStartPeriod=60s
HealthStartupInterval=5s
NoNewPrivileges=true
Secret=db_password
PodmanArgs=--ulimit=nofile=65536:65536
NetworkAlias=postgres
Notify=healthy
Label=io.podman.quadlet.project=example-proxy
Label=io.podman.quadlet.service=postgres
AutoUpdate=registry

[Service]
CPUQuota=200%
AllowedCPUs=2,3,4,5
MemoryMax=1g
MemoryLow=512m
MemorySwapMax=1g
TasksMax=256
OOMScoreAdjust=-900
ExecStopPre=podman exec postgres su-exec postgres pg_ctl stop -m fast
Restart=always

[Install]
WantedBy=default.target
