# data/examples/48/actual.container
# nginx.container
[Container]
PublishPort=80:80
PublishPort=443:443
Volume=wellknown:/usr/share/nginx/html/.well-known

[Service]
Restart=always
