# data/fixtures/env-ports/output/metrics.container
[Container]
Image=prom/prometheus:latest
Pod=env-ports.pod
