# data/fixtures/env-ports/output/redis.container
[Container]
Image=redis:7-alpine
Pod=env-ports.pod
