# data/examples/45/actual.container
### webapiapplication.container ###
[Container]
Image=webapiapplication
