# data/examples/45/example.container
# webapiapplication.container
[Container]
Image=webapiapplication
Label=kompose.service.type=NodePort
