# data/examples/15/actual.container
# app.container
[Container]
PublishPort=3019:3019
Environment=DATABASE_URL=postgres://app:app@db:5432/app
