# data/examples/24/actual.container
# note-directory-classifier.container
[Container]
Volume=./data:/data
