# data/examples/39/actual.container
# db.container
[Container]
Image=postgres:10.1
PublishPort=5432
Volume=db-data:/var/lib/postgresql/data
