# data/examples/118/actual.container
### jupyter.container ###
[Container]
Image=continuumio/miniconda3
PublishPort=8998:8888

[Service]
Restart=always
