# data/examples/36/actual.container
# web.container
[Container]
Image=centos/httpd
PublishPort=3030:3000
Volume=./app:/src/app
