# data/examples/117/actual.container
# app.container
[Container]
Image=eibrahim/fluid-calendar:latest
PublishPort=${CAL_PORT}:3000

[Service]
Restart=always
