# data/examples/93/actual.container
# katex.container
[Container]
PublishPort=57359:8080
ContainerName=katex
