# data/examples/82/actual.container
### lilium.container ###
[Container]
PublishPort=${PORT}:${PORT}
Volume=./.data:/app/.data:Z
ContainerName=lilium

[Service]
Restart=always
