# data/examples/98/actual.container
# postgres.container
[Container]
Image=postgres:16-alpine
Volume=./data/pgdata:/var/lib/postgresql/data

[Service]
Restart=always
