# data/examples/83/actual.container
### tabby.container ###
[Container]
Image=tabbyml/tabby
PublishPort=7652:8080
Volume=$HOME/.tabby:/data

[Service]
Restart=always
