# data/examples/88/actual.container
# suse.container
[Container]
Image=opensuse-datascience
Volume=/tmp/:/tmp/
Environment=DISPLAY:{$DISPLAY}
