# data/examples/41/actual.container
### client.container ###
[Container]
Image=registry.centos.org/centos/centos:7
PublishPort=1337
