# data/examples/38/actual.container
### httpd.container ###
[Container]
Image=docker.io/fedora/apache
PublishPort=80
Volume=/var/www/html
