# data/fixtures/env-ports/output/env-ports/db.container
[Container]
Image=postgres:17-alpine
Pod=env-ports.pod
Environment=POSTGRES_PASSWORD=secret
