# data/fixtures/env-ports/output/env-ports/web.container
[Container]
Image=nginx:alpine
Pod=env-ports.pod
Environment=BACKEND_URL=http://app:${APP_PORT}
