# data/examples/104/actual.container
# app.container
[Container]
Image=eibrahim/fluid-calendar:latest
PublishPort=3000:3000

[Service]
Restart=always
