# data/examples/57/actual.container
# db.container
[Container]
PublishPort=32789:5432
Volume=pgdata:/var/lib/postgresql/data
