# data/examples/6/actual.container
# dokuwiki.container
[Container]
Image=dokuwiki/dokuwiki:stable
Volume=./storage:/storage
ContainerName=dokuwiki
