# data/examples/6/example.container
# dokuwiki.container
[Container]
ContainerName=dokuwiki
Image=dokuwiki/dokuwiki:stable
NoNewPrivileges=true
Volume=./storage:/storage
