# data/examples/60/actual.container
### nfs.container ###
[Container]
PublishPort=2049:2049
Volume=./data:/exports
