# data/examples/42/example.container
# client.container
[Container]
Image=registry.centos.org/centos/centos:7
PodmanArgs=--attach stdin
PublishPort=1337
