# data/examples/84/example.container
# jsonhero_web.container
[Container]
ContainerName=super-productivity
Image=johannesjo/super-productivity
PublishPort=8782:80
