# data/examples/40/actual.container
### db.container ###
[Container]
Image=postgres:10.1
PublishPort=5432
Volume=./data_sux:/var/lib/postgresql/data_sux
