# data/examples/37/actual.container
# redis-master.container
[Container]
Image=registry.k8s.io/redis:e2e
PublishPort=6379
